# files/supervisord_base.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700                       ; sockef file mode (default 0700)

[supervisord]
minfds=10240; (min. avail startup file descriptors;default 1024)
stdout_logfile=/proc/self/fd/1
stderr_logfile=/proc/self/fd/2
logfile=/proc/self/fd/2
pidfile=/tmp/supervisord.pid  ; (supervisord pidfile;default supervisord.pid)
childlogdir=/var/log/supervisor            ; ('AUTO' child log dir, default $TEMP)
logfile_maxbytes=0         ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=0            ; (num of main logfile rotation backups;default 10)
loglevel=debug                 ; (log level;default info; others: debug,warn,trace)
nodaemon=true                 ; (start in foreground if true;default false)
minfds=1024                   ; (min. avail startup file descriptors;default 1024)
minprocs=200                  ; (min. avail process descriptors;default 200)
user=root                     ;

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:syslogd]
command=/usr/sbin/rsyslogd -n
autostart=true
autorestart=true
startretries=20
priority=5
stdout_events_enabled=true
stdout_logfile=/proc/self/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0

[program:php-fpm]
command=/bin/run_with_env.sh /site/web/.envDocker /usr/sbin/php-fpmPHP_VERSION --nodaemonize --force-stderr -c /etc/php/PHP_VERSION/fpm --fpm-config /etc/php/PHP_VERSION/fpm/php-fpm.conf
autostart=true
autorestart=true
startretries=20
priority=10
user=root
numprocs=ENABLE_WEB
stdout_events_enabled=true
stdout_logfile=/proc/self/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0

[program:cron]
command=/usr/sbin/cron -f
autostart=true
autorestart=true
startretries=20
priority=20
stdout_events_enabled=true
stdout_logfile=/proc/self/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -c /site/nginx/config/nginx.conf
autostart=true
autorestart=true
startretries=20
priority=30
user=root
numprocs=ENABLE_WEB
stdout_events_enabled=true
stdout_logfile=/proc/self/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0

[program:laravel-horizon]
command=/bin/run_with_env.sh /site/web/.envDocker /usr/bin/php artisan horizon
autostart=true
autorestart=true
startretries=20
priority=25
process_name=%(program_name)s_%(process_num)02d
directory=/site/web/
user=web
numprocs=ENABLE_HORIZON
stdout_events_enabled=true
stdout_logfile=/proc/self/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0
; 3 hours
stopwaitsecs=10800

[program:laravel-pulse-check]
command=/bin/run_with_env.sh /site/web/.envDocker /usr/bin/php artisan pulse:check
autostart=true
autorestart=true
startretries=20
priority=25
process_name=%(program_name)s_%(process_num)02d
directory=/site/web/
user=web
numprocs=ENABLE_PULSE_CHECK
stdout_events_enabled=true
stdout_logfile=/proc/self/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0
; 3 hours
stopwaitsecs=10800

[program:laravel-pulse-work]
command=/bin/run_with_env.sh /site/web/.envDocker /usr/bin/php artisan pulse:work
autostart=true
autorestart=true
startretries=20
priority=25
process_name=%(program_name)s_%(process_num)02d
directory=/site/web/
user=web
numprocs=ENABLE_PULSE_WORK
stdout_events_enabled=true
stdout_logfile=/proc/self/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0
; 3 hours
stopwaitsecs=10800

[program:laravel-websocket]
command=/bin/run_with_env.sh /site/web/.envDocker /usr/bin/php artisan websockets:serve --port=LARAVEL_WEBSOCKETS_PORT
priority=25
process_name=%(program_name)s_%(process_num)02d
startretries=20
autostart=true
autorestart=true
user=web
numprocs=ENABLE_WEBSOCKET
stdout_events_enabled=true
stdout_logfile=/proc/self/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0
; 3 hours
stopwaitsecs=10800

[program:laravel-worker]
priority=25
process_name=%(program_name)s_%(process_num)02d
command=/bin/run_with_env.sh /site/web/.envDocker /usr/bin/php artisan queue:work --sleep=3 --tries=3 --max-time=3600
startretries=20
autostart=true
autorestart=true
user=web
numprocs=SIMPLE_WORKER_NUM
stdout_events_enabled=true
stdout_logfile=/proc/self/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0
; 3 hours
stopwaitsecs=10800


[program:sshd]
priority=30
process_name=%(program_name)s_%(process_num)02d
command=/usr/sbin/sshd -D
startretries=20
autostart=true
autorestart=true
user=root
numprocs=ENABLE_SSH
stdout_events_enabled=true
stdout_logfile=/proc/self/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0
; 3 hours
stopwaitsecs=10800
