# files/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700                       ; sockef file mode (default 0700)

[supervisord]
logfile=/site/logs/supervisor/supervisord.log  ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB         ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10            ; (num of main logfile rotation backups;default 10)
loglevel=info                 ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid  ; (supervisord pidfile;default supervisord.pid)
nodaemon=true                 ; (start in foreground if true;default false)
minfds=1024                   ; (min. avail startup file descriptors;default 1024)
minprocs=200                  ; (min. avail process descriptors;default 200)
user=root                     ;

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:syslogd]
priority=5
autostart=true
autorestart=true
command=/usr/sbin/rsyslogd -n
stdout_events_enabled=true
; stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_logfile=/site/logs/supervisor/%(program_name)s.log
; stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stderr_logfile=/site/logs/supervisor/%(program_name)s.log

[program:php-fpm]
command=/usr/sbin/php-fpmPHP_VERSION --nodaemonize -c /etc/php/PHP_VERSION/fpm --fpm-config /etc/php/PHP_VERSION/fpm/php-fpm.conf
autostart=true
autorestart=true
priority=10
stdout_events_enabled=true
; stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_logfile=/site/logs/supervisor/%(program_name)s.log
; stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stderr_logfile=/site/logs/supervisor/%(program_name)s.log

[program:cron]
command=/usr/sbin/cron -f
autostart=true
autorestart=true
priority=20
stdout_events_enabled=true
; stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_logfile=/site/logs/supervisor/%(program_name)s.log
; stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stderr_logfile=/site/logs/supervisor/%(program_name)s.log

[program:laravel-horizon]
priority=25
process_name=%(program_name)s_%(process_num)02d
directory=/site/web/
command=/bin/run_with_env.sh /site/web/.env /usr/bin/php artisan horizon
autostart=true
autorestart=true
user=web
numprocs=ENABLE_HORIZON
stdout_events_enabled=true
; stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_logfile=/site/logs/supervisor/%(program_name)s.log
; stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stderr_logfile=/site/logs/supervisor/%(program_name)s.log

[program:nginx]
command=/usr/sbin/nginx -c /site/nginx/config/nginx.conf
autostart=true
autorestart=true
priority=30
stdout_events_enabled=true
; stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_logfile=/site/logs/supervisor/%(program_name)s.log
; stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stderr_logfile=/site/logs/supervisor/%(program_name)s.log

[program:filebeat]
command=bash -c "sleep 15 && /usr/share/filebeat/bin/filebeat -path.home /usr/share/filebeat -path.config /etc/filebeat -path.data /var/lib/filebeat -path.logs /var/log/filebeat"
autostart=ELK_FILEBEAT_ACTIVE
autorestart=ELK_FILEBEAT_ACTIVE
priority=30
stdout_events_enabled=true
; stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_logfile=/site/logs/supervisor/%(program_name)s.log
; stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stderr_logfile=/site/logs/supervisor/%(program_name)s.log

[program:metricbeat]
command=bash -c "sleep 15 && /usr/share/metricbeat/bin/metricbeat -path.home /usr/share/metricbeat -path.config /etc/metricbeat -path.data /var/lib/metricbeat -path.logs /var/log/metricbeat"
autostart=ELK_METRICBEAT_ACTIVE
autorestart=ELK_METRICBEAT_ACTIVE
priority=30
stdout_events_enabled=true
; stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_logfile=/site/logs/supervisor/%(program_name)s.log
; stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stderr_logfile=/site/logs/supervisor/%(program_name)s.log
